# operator_console/systemd/launch-manager.service
[Unit]
Description=GO2W Navigation Launch Manager
Documentation=https://github.com/your-org/go2w_real
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=unitree
Group=unitree
WorkingDirectory=/home/unitree/ros_ws/src/go2w_real/operator_console
ExecStart=/usr/bin/python3 /home/unitree/ros_ws/src/go2w_real/operator_console/launch_manager.py --port 8080
Restart=always
RestartSec=5
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=20
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
